# hyperiond
#!/sbin/runscript
# Hyperion daemon
# description: Hyperion daemon
# processname: hyperiond

depend() {
	after net
}

start() {
	ebegin "Starting hyperiond server..."
	start-stop-daemon --start --exec /usr/bin/hyperiond /etc/hyperion.config.json --pidfile /var/run/hyperiond.pid -b -m
	eend $?
}

stop() {
	ebegin "Stopping hyperiond server..."
	start-stop-daemon --stop --exec /usr/bin/hyperiond stop --pidfile /var/run/hyperiond.pid
	eend $?
}
